# launchd/com.mama.whatsapp.tick.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.mama.whatsapp.tick</string>

    <key>ProgramArguments</key>
    <array>
      <string>/Users/alexbuchel/projects/mama_whatsapp/launchd/run-tick.sh</string>
    </array>

    <key>StartInterval</key>
    <integer>3600</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/mama_whatsapp_tick.out.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/mama_whatsapp_tick.err.log</string>
  </dict>
</plist>
